# systemd/fido-ssh-agent.service
[Unit]
Description=FIDO SSH Agent
Requires=fido-ssh-agent.socket

[Service]
ExecStart=%h/.cargo/bin/fido-ssh-agent
Restart=on-failure
