# systemd/fido-ssh-agent.socket
[Unit]
Description=FIDO SSH Agent socket

[Socket]
ListenStream=%t/fido-ssh-agent.sock
SocketMode=0600

[Install]
WantedBy=sockets.target
