# systemd/fido-ssh-agent.socket
[Unit]
Description=FIDO SSH Agent socket
Before=graphical-session-pre.target

[Socket]
ListenStream=%t/fido-ssh-agent.sock
ExecStartPost=/bin/systemctl --user set-environment SSH_AUTH_SOCK=%t/fido-ssh-agent.sock

[Install]
WantedBy=sockets.target
